# etc/supervisord.d/fcgi.ini
[program:spawn-fcgi]
process_name = fcgiwrap
command=/etc/init.d/spawn-fcgi start
startsecs=0
autorestart=false
