# etc/supervisord.d/mailman.ini
[program:mailman]
process_name = qrunner
command=/usr/share/mailman/bin/mailmanctl -q start
startsecs=0
autorestart=false
priority=90