# etc/supervisord.d/postfix.ini
[program:postfix]
process_name = master
command=/etc/init.d/postfix start
startsecs=0
autorestart=false
priority=99