# etc/supervisord.d/rsyslog.ini
[program:rsyslog]
priority=0
command=/sbin/rsyslogd -n
